# deploy/prata-backend.service
[Unit]
Description=prata-sys FastAPI backend
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=prata
Group=prata
WorkingDirectory=/opt/prata-sys/backend
ExecStart=/opt/prata-sys/backend/.venv/bin/python -m uvicorn app.main:app \
    --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/prata-sys/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/prata-sys/backend
PrivateTmp=true

[Install]
WantedBy=multi-user.target
